# deploy/systemd/ai-manager-backend.service
[Unit]
Description=AI Manager Admin Backend
After=network.target

[Service]
Type=simple
User=aimanager
Group=aimanager
WorkingDirectory=/opt/ai-manager/backend
Environment=TZ=Asia/Shanghai
EnvironmentFile=/opt/ai-manager/backend/backend.env
ExecStart=/usr/bin/java -Duser.timezone=Asia/Shanghai -Xms256m -Xmx512m -jar /opt/ai-manager/backend/admin-server.jar
Restart=on-failure
RestartSec=10

# 本机一键部署需 aimanager 调用 sudo（setuid），不可开启 NoNewPrivileges
NoNewPrivileges=no

# 上传目录与日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-manager-backend

[Install]
WantedBy=multi-user.target
